# root/etc/supervisor/conf.d/07_greyhole-boot-init.conf
[program:greyhole-boot-init]
priority=101
command=/usr/bin/greyhole --boot-init
user=root
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=false
startretries=0
startsecs=0
